# ssh/rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

/usr/sbin/sshd -f /etc/ssh/sshd_config -D
