# ssh/rootfs/etc/services.d/ssh/run
#!/bin/sh

ssh-keygen -A
exec /usr/sbin/sshd -f /etc/ssh/sshd_config -D
