# base/rootfs/etc/services.d/syslog-ng/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/syslog-ng -F -v --pidfile /run/syslog-ng.pid
