# base/rootfs/etc/services.d/syslog-ng/run
#!/bin/sh

exec /usr/sbin/syslog-ng -F --pidfile /run/syslog-ng.pid
