# base/rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/crond -f -S -l 0 -d 8 -c /etc/crontabs
