# base/rootfs/etc/services.d/cron/run
#!/bin/execlineb -P

exec -c
fdmove -c 2 1

/usr/sbin/crond -f -d 0 -c /etc/crontabs
