# mariadb/rootfs/etc/services.d/mariadb/run
#!/usr/bin/with-contenv sh

exec /usr/bin/mysqld_safe --syslog --pid-file=/var/run/mysqld/mysqld.pid
